# b2/rootfs/etc/services.d/backblaze/run
#!/usr/bin/with-contenv bashio
BUCKET=$(bashio::config 'bucket')
declare -a options

if bashio::var.has_value "$(bashio::config 'keep_days')"; then
  options+=(--keepDays "$(bashio::config 'keep_days')")
else
  options+=(--delete)
fi

if bashio::config.true 'allowEmptySource'; then
  options+=(--allowEmptySource)
fi

if [[ $(bashio::config 'log_level') != 'debug' ]]; then
  options+=(--noProgress)
fi

b2 sync "${options[@]}" /backup "b2://$BUCKET"
