# ngrok/rootfs/etc/services.d/ngrok/run
#!/usr/bin/with-contenv bashio

exec /usr/bin/ngrok.sh
