# cloudflared/rootfs/etc/services.d/cloudflared/finish
#!/usr/bin/env bashio

/run/s6/basedir/bin/halt