# tunnels/rootfs/etc/services.d/tunnels/run
#!/usr/bin/with-contenv bashio
declare -a options
declare attempts

attempts=3

options+=" -p $(bashio::config 'port')"
options+=" -i /share/ssh-tunnels/id_rsa"
options+=" -o StrictHostKeyChecking=no"

for id in $(bashio::config "tunnels|keys"); do
  tunnel=$(bashio::config "tunnels[${id}].tunnel")
  bashio::log.debug "Adding tunnel to request: $tunnel"
  options+=" -R $tunnel"
done

SERVER=$(bashio::config 'server')
if [[ $(bashio::config.has_value 'user') ]]; then
  USER=$(bashio::config 'user')
  options+=" $USER@$SERVER"
else
  options+=" $SERVER"
fi

while [[ $attempts -gt 0 ]]; do
  let "attempts-=1"
  bashio::log.debug "Executing: ssh ${options[@]}"
  ssh "${options[@]}"
  bashio::log.debug "End SSH session..."
  bashio::log.warning "Failed to establish tunnel"
  bashio::log.info "Waiting 30 seconds before trying again"
  sleep 30
done

bashio::log.error "Failed to establish a tunnel too many times, please check your settings and try again."
exit 1
